# supervisor/cacherefresher.ini
[program:cacherefresher]

command=/usr/local/bin/python /code/cacherefresher.py
autostart=true
autorestart=true